# query-pdb/supervisord.conf
[supervisord]
nodaemon=true

[program:query-pdb-server]
command=/query-pdb/build/server/query_pdb_server
autostart=true
autorestart=true
